# docker-kinc/images/executor/supervisord/autofs.ini
[program:autofs]
command=/usr/sbin/automount --pid-file /run/autofs.pid
autostart=true
autorestart=true
